# services/hics-web.service
[Unit]
Description=HICS Web Dashboard (Flask, port 5000)
After=multi-user.target hics-core.service

[Service]
Type=simple
User=pawan
WorkingDirectory=/home/pawan/hics
ExecStart=/usr/bin/python3 /home/pawan/hics/web/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
